# deploy/systemd/diary-backup.service
# deploy/systemd/diary-backup.service
[Unit]
Description=diary nightly SQLite backup
After=diary.service
Wants=diary.service
[Service]
Type=oneshot
ExecStart=%h/work/website/diary/deploy/scripts/diary-backup.sh
