# deploy/systemd/unflincher-backup.service
# deploy/systemd/unflincher-backup.service
[Unit]
Description=unflincher nightly SQLite backup
After=unflincher.service
Wants=unflincher.service
[Service]
Type=oneshot
# unflincher-backup.sh and verify-unflincher-backup.py are installed together under
# %h/.local/bin/ (see README's "Backups & recovery" section). The backup script resolves the
# verifier beside itself, so re-copy both files after pulling a change to either one.
ExecStart=%h/.local/bin/unflincher-backup.sh
