# deploy/quadlet/diary.container
# deploy/quadlet/diary.container
[Unit]
Description=diary.yourdomain.com — AI-annotated diary (FastAPI + htmx + SSE, in-process async batch worker)

[Container]
Image=localhost/diary:latest
ContainerName=diary

# Loopback only. cloudflared maps diary.yourdomain.com -> http://localhost:8096 (see
# deploy/cloudflared/diary-ingress.snippet.yml, Task 20). PRIVATE behind Cloudflare Access.
PublishPort=127.0.0.1:8096:8000

Volume=diary-data.volume:/data

# Shared secret — the SAME fine-grained GitHub PAT trading's and News Center's Copilot SDK
# wrappers use (same GitHub account, same Copilot subscription). See website/README.md's secrets
# section; don't fork a diary-only copy. CopilotClient() auto-detects this env var.
Secret=unflincher-copilot-github-token,type=env,target=COPILOT_GITHUB_TOKEN

Environment=DIARY_DB=/data/diary.db
Environment=DIARY_LLM_MODEL=claude-sonnet-4.6
Environment=DIARY_BATCH_CONCURRENCY=3
# CF_TEAM_DOMAIN and OPERATOR_EMAIL match shunlyu-infra/deploy/infra.env's shared, non-secret
# account facts (same values every Access-gated site on this box uses). AUD is the one value
# that's genuinely per-app -- only known after running deploy/create-access-diary-app.sh below.
Environment=DIARY_CF_TEAM_DOMAIN=your-team-name.cloudflareaccess.com
Environment=DIARY_CF_ACCESS_AUD=REPLACE_WITH_YOUR_ACCESS_APP_AUD
Environment=DIARY_OPERATOR_EMAIL=your-email@example.com
Environment=DIARY_REQUIRE_ACCESS_AUTH=true

HealthCmd=python -c "import urllib.request,sys; sys.exit(0 if urllib.request.urlopen('http://127.0.0.1:8000/healthz',timeout=5).status==200 else 1)"
HealthInterval=60s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=30s
# Deliberately NO HealthOnFailure=kill — a health blip must never interrupt an in-flight
# apply-to-all batch (technical design §7.4/§7.6). Restart=always below covers real crashes.

PodmanArgs=--stop-timeout=30

[Service]
Restart=always
RestartSec=5
TimeoutStartSec=90
TimeoutStopSec=45
# uvicorn exits 143 (128+SIGTERM) on a clean graceful shutdown -- without this, systemd logs
# every intentional stop/restart as "Failed with result 'exit-code'" even though the app itself
# shut down cleanly within TimeoutStopSec (verified: real shutdown took 0.3s).
SuccessExitStatus=143

[Install]
WantedBy=default.target
