# deploy/quadlet/unflincher.container
# deploy/quadlet/unflincher.container
[Unit]
Description=unflincher.yourdomain.com — AI-annotated diary (FastAPI + htmx + SSE, in-process async batch worker)

[Container]
Image=localhost/unflincher:latest
ContainerName=unflincher

# Loopback only. cloudflared maps unflincher.yourdomain.com -> http://localhost:8096 (see
# deploy/cloudflared/unflincher-ingress.snippet.yml). PRIVATE behind Cloudflare Access.
PublishPort=127.0.0.1:8096:8000

Volume=unflincher-data.volume:/data

# If you already run other apps on this host that use the GitHub Copilot SDK, they can all
# share ONE podman secret (same GitHub account, same Copilot subscription) rather than each
# minting their own. CopilotClient() auto-detects this env var. If you rename this secret,
# also update UNFLINCHER_COPILOT_SECRET when running deploy/scripts/deploy-unflincher.sh so
# the two stay in sync.
#
# NOTE: this project was renamed from "diary" to "unflincher"; the podman secret name below
# deliberately was NOT renamed along with everything else, because podman secrets are
# write-only (there is no way to read the token back out to recreate it under a new name
# without the owner re-supplying the raw token). "diary-copilot-github-token" is the correct,
# intentional value here even in an otherwise fully-renamed project.
Secret=diary-copilot-github-token,type=env,target=COPILOT_GITHUB_TOKEN

Environment=UNFLINCHER_DB=/data/unflincher.db
Environment=UNFLINCHER_LLM_MODEL=claude-sonnet-4.6
Environment=UNFLINCHER_BATCH_CONCURRENCY=3
Environment=UNFLINCHER_LLM_CONCURRENCY=4
# CF_TEAM_DOMAIN is the SHORT team name only (this app's own auth.py appends
# ".cloudflareaccess.com" itself for both the JWKS URL and the issuer check) -- if you're
# copying this value from Cloudflare's own dashboard/JWKS URL, that shows the full
# "<team>.cloudflareaccess.com" form instead, so strip the suffix before pasting it here.
# AUD is the one value that's genuinely per-app -- only known after running
# deploy/create-access-unflincher-app.sh below.
Environment=UNFLINCHER_CF_TEAM_DOMAIN=your-team-name
Environment=UNFLINCHER_CF_ACCESS_AUD=REPLACE_WITH_YOUR_ACCESS_APP_AUD
Environment=UNFLINCHER_OPERATOR_EMAIL=your-email@example.com
Environment=UNFLINCHER_REQUIRE_ACCESS_AUTH=true

HealthCmd=python -c "import urllib.request,sys; sys.exit(0 if urllib.request.urlopen('http://127.0.0.1:8000/healthz',timeout=5).status==200 else 1)"
HealthInterval=60s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=30s
# Deliberately NO HealthOnFailure=kill — a health blip must never interrupt an in-flight
# apply-to-all batch job. Restart=always below covers real crashes.

PodmanArgs=--stop-timeout=30

[Service]
Restart=always
RestartSec=5
TimeoutStartSec=90
TimeoutStopSec=45
# uvicorn exits 143 (128+SIGTERM) on a clean graceful shutdown -- without this, systemd logs
# every intentional stop/restart as "Failed with result 'exit-code'" even though the app itself
# shut down cleanly within TimeoutStopSec (verified: real shutdown took 0.3s).
SuccessExitStatus=143

[Install]
WantedBy=default.target
